# HomeAssistantWireCat/rootfs/etc/s6-overlay/s6-rc.d/wireguard/finish
#!/command/with-contenv bashio
# Tear the interface down whenever the wireguard service stops or restarts.
wg-quick down wg0 2>/dev/null || true
